# templates/mysql_purge_logs.timer.j2
[Unit]
Description=MySQL Binary Log Purge Timer

[Timer]
OnCalendar=*-*-* *:0,15,30,45:00
Unit=mysql_purge_logs.service

[Install]
WantedBy=timers.target
